# packages/mediacenter/emby-dependencies/kodi-theme-Embuary/system.d/kodi-theme-Embuary.service
[Unit]
Description=Embuary addons copy script
Before=kodi.service
After=network-online.target graphical.target

[Service]
Type=oneshot
Environment=HOME=/storage
ExecStart=/usr/bin/kodi-theme-Embuary
RemainAfterExit=yes

[Install]
WantedBy=kodi.service